# pkg/debiansync/Library/LaunchDaemons/com.debiansync.fullsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.debiansync.fullsync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/debian-sync</string>
    <string>full</string>
    <string>--quiet</string>
  </array>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/var/log/debiansync.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/debiansync.log</string>
</dict>
</plist>
